# service/timer/paccache.timer
[Unit]
Description = clearing pacman cache in timer

[Timer]
OnCalendar = week
Persistent = true

[Install]
WantedBy = multi-user.target
